# systemd/user/kdeinit.service
[Unit]
Description=KDE 'kdeinit' - TODO

[Service]
ExecStart=/usr/bin/kdeinit5 +kcminit_startup --no-fork
KillMode=none
BusName=org.kde.klauncher5
