# systemd/user/kf5-plasma.target
[Unit]
Description=KDE Frameworks 5/Plasma
Wants=wm.target

[Install]
Alias=default.target
