# systemd/user/ksmserver.service
[Unit]
Description=KDE Session Management Server

[Service]
ExecStart=/usr/bin/ksmserver
BusName=org.kde.ksmserver
