# systemd/user/ksmserver.service
[Unit]
Description=KDE Session Management Server
Requires=kglobalaccel.service

[Service]
ExecStart=/usr/bin/ksmserver
BusName=org.kde.ksmserver
