# systemd/user/plasma-workspace.target
[Unit]
Description=KDE Plasma Workspace
Wants=plasma-workspace.service window-manager.service launcher.service
